# templates/confluence.service.erb
[Unit]
Description=Confluence Team Collaboration Software
After=syslog.target network.target

[Service]
Type=forking
User=<%= scope.lookupvar('confluence::user') %>
Environment="JAVA_HOME=<%= scope.lookupvar('confluence::javahome') %>"
PIDFile=<%= scope.lookupvar('confluence::webappdir') %>/work/catalina.pid
ExecStart=<%= scope.lookupvar('confluence::webappdir') %>/bin/start-confluence.sh
ExecStop=<%= scope.lookupvar('confluence::webappdir') %>/bin/stop-confluence.sh

[Install]
WantedBy=multi-user.target
